# crontab
* * * * * /bin/rho/sync_themes.sh >> /var/log/cron.log 2>&1
